# jenkins_scripts/bedade_web_live.conf
[program:bedade_web]
command = /webapps/bedade/bedade/bedade_web.sh
user = bedade
stdout_logfile = /webapps/bedade/bedade/supervisor.log
redirect_stderr = true