# jenkins_scripts/bedade_web.conf
[program:bedade_web]
command = /root/projects/bedade-internal/bedade_web_dev.sh
user = root
stdout_logfile = /root/projects/bedade-internal/supervisor.log
redirect_stderr = true
killasgroup=true