# jenkins_scripts/bedade_celeryworker_live.conf
[program:bedade_celeryworker]
command = /webapps/bedade/bedade/bedade_celeryworker.sh
user = bedade
stdout_logfile = /webapps/bedade/bedade/supervisor.log
redirect_stderr = true