# jenkins_scripts/bedade_celeryworker.conf
[program:bedade_celeryworker]
command = /root/projects/bedade-internal/bedade_celeryworker_dev.sh
user = root
stdout_logfile = /root/projects/bedade-internal/supervisor.log
redirect_stderr = true